# etc/sv/syslogd/run
#!/bin/sh
. /etc/default/syslogd
umask "$UMASK"
exec syslogd -n $OPTS
