# etc/sv/agetty-tty6/run
#!/bin/sh
exec chpst -P /sbin/agetty tty6 38400 linux
