# etc/sv/agetty-tty1/run
#!/bin/sh
exec chpst -P /sbin/agetty tty1 38400 linux
