# etc/sv/dbus/run
#!/bin/sh
. /etc/default/dbus
exec dbus-daemon --nofork $OPTS
